# infra/cloudbeds-enrichment.service
[Unit]
Description=Cloudbeds Enrichment Consumer
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sadie_gtm_ingestor
ExecStart=/home/ubuntu/.local/bin/uv run python -m workflows.enrich_cloudbeds_consumer --concurrency 10
Restart=always
RestartSec=10
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
